# com.nerdlab.iboss-qna-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nerdlab.iboss-qna-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jaewansim/Documents/nerdlab/news-crawler-bot/src/iboss_qna_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jaewansim/Documents/nerdlab/news-crawler-bot</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jaewansim/Documents/nerdlab/news-crawler-bot/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jaewansim/Documents/nerdlab/news-crawler-bot/logs/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>